# themes/luci-theme-dark/root/etc/uci-defaults/30_luci-theme-dark
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Dark=/luci-static/dark
	set luci.main.mediaurlbase=/luci-static/dark
	commit luci
EOF
exit 0
