# themes/luci-theme-carbone/root/etc/uci-defaults/30_luci-theme-carbone
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Dark=/luci-static/carbone
	set luci.main.mediaurlbase=/luci-static/carbone
	commit luci
EOF
exit 0
